# roles/firewall/templates/miuops-firewall-docker-user.service.j2
[Unit]
Description=miuops container isolation (DOCKER-USER) — applied with Docker + on a timer
# After+PartOf+WantedBy=docker.service, but NO Requires: start with/after Docker and
# re-run on its restart, WITHOUT pulling Docker back up during maintenance or recording
# a false failure before Docker is installed. (The applier skips its FORWARD-jump assert
# when Docker isn't active, so a timer fire during a Docker stop is a clean no-op.)
After=docker.service
PartOf=docker.service

[Service]
Type=oneshot
# No RemainAfterExit: the reconcile .timer re-runs ExecStart on each fire. (A
# RemainAfterExit oneshot stays "active (exited)", so the timer's `start` would be a
# no-op and the periodic self-heal would silently never run.)
# TimeoutStartSec caps a pathological hang (e.g. a wedged xtables lock) so the unit fails
# loudly rather than staying "activating"; normal runs finish in well under a second.
TimeoutStartSec=180
ExecStart=/usr/local/sbin/miuops-firewall.sh docker-user

[Install]
# WantedBy docker.service (not multi-user.target): `systemctl enable` puts this in
# docker.service.wants/, so it runs right after Docker starts (first install + every
# boot), applying DOCKER-USER before any role publishes a container port.
WantedBy=docker.service
