# roles/firewall/templates/miuops-firewall-docker-user.timer.j2
[Unit]
Description=miuops periodic DOCKER-USER reconcile (self-heal drift)
[Timer]
OnBootSec={{ firewall_reconcile_interval }}
OnUnitActiveSec={{ firewall_reconcile_interval }}
[Install]
WantedBy=timers.target
