# roles/firewall/templates/miuops-firewall-docker-user.service.j2
[Unit]
Description=miuops container isolation (DOCKER-USER), re-applied after Docker
After=docker.service
Requires=docker.service
PartOf=docker.service
[Service]
Type=oneshot
RemainAfterExit=yes
# Cap a pathological hang (e.g. a wedged xtables lock) so the unit fails loudly
# rather than staying "activating" forever; normal runs complete in well under a second.
TimeoutStartSec=180
ExecStart=/usr/local/sbin/miuops-firewall.sh docker-user
[Install]
WantedBy=multi-user.target
